# PVBattery
#!/sbin/openrc-run
# Distributed under the terms of the GNU General Public License v2
# $Header: $

description="Start control of a PV Battery"
command="/usr/bin/lua"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="/opt/PVBattery/PVBatteryGuard.lua"

command_background=true
command_user="root:root"

directory="/opt/PVBattery/"

rm -f "/tmp/index.html"
ln -sf "/tmp/index.html" "/var/www/localhost/htdocs/index.html"

date > /tmp/index.html
echo -e "\nPVBattery is starting up ..." >> /tmp/index.html
sync
